# runit/uplink/run
#!/bin/sh

cd /usr/share/bytebeam/uplink

exec ./target/debug/uplink -c configs/simulator.toml -a configs/simulator.json -vv
