# runit/uplink/run
#!/bin/sh

cd /usr/share/bytebeam/uplink

export RUST_LOG=debug
exec ./target/debug/uplink -c configs/demo.toml -a configs/noauth.json
