# runit/uplink/run
#!/bin/sh

cd /usr/share/bytebeam/uplink

exec uplink -a /usr/share/bytebeam/uplink/shared/device.json
