# runit/simulator/finish
#!/bin/sh

cd /usr/share/bytebeam/uplink

exec ./simulator.sh kill_devices
echo "Killed all simulated devices"
